# discovery-service/applypilot-discovery.service
[Unit]
Description=ApplyPilot Discovery Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=applypilot
WorkingDirectory=/opt/applypilot/discovery-service
ExecStart=/usr/bin/env python3.11 main.py
Restart=on-failure
RestartSec=30

# Load env from the service's .env file
EnvironmentFile=/opt/applypilot/discovery-service/.env

# Stdout → journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=applypilot-discovery

[Install]
WantedBy=multi-user.target
